# sarissa.service
[Unit]
Description=Sarissa Security Platform
After=network.target postgresql.service

[Service]
Type=simple
User=sarissa
WorkingDirectory=/opt/sarissa
ExecStart=/opt/sarissa/sarissa
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
StandardOutput=append:/var/log/sarissa/sarissa.log
StandardError=append:/var/log/sarissa/sarissa.err
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target 